# init/zsh/zshrc.d/local-functions/etc/___TODOADDITION___/CORE/LINUX/kwtools-master/etc/skel/.config/autostart/firstlogin.desktop

[Desktop Entry]
Encoding=UTF-8
Version=1.0
Type=Application
Terminal=false
Name[de_DE]=firstlogin
Exec=/usr/bin/firstlogin
Name=firstlogin
